# mongroup.conf
description 'mongroup start'
author 'maxogden'
start on runlevel [0123456]
task
exec start-stop-daemon --start -u $USER --chdir $USER_HOME --exec /usr/bin/env mongroup start
